# packaging/systemd/voice-dictation-gui.service
[Unit]
Description=Voice Dictation GUI Overlay
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/voice-dictation gui
ExecStop=/usr/bin/pkill -TERM -f "voice-dictation gui"
Restart=no
Environment="RUST_LOG=info"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
