# packaging/systemd/voice-dictation.service
[Unit]
Description=Voice Dictation Persistent Daemon with Integrated GUI
After=pipewire.service graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/voice-dictation daemon
Restart=on-failure
RestartSec=5
Environment="RUST_LOG=info"
Environment="GUI_LOG=info"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
